# configs/releng_orchid/airootfs/etc/systemd/system/multi-user.target.wants/first-boot.service
[Unit]
Description=Orchid First Boot

[Service]
Type=oneshot
ExecStart=/etc/orchid/first_boot.sh

[Install]
WantedBy=multi-user.target
